# files/etc/init.d/acctl
#!/bin/sh /etc/rc.common
#
# AC Controller OpenWrt Init Script (procd)
# /etc/init.d/acctl — starts AC server (acser)
#
# This script is OpenWrt-standard:
#   start_service() -> starts service via procd (respawn support)
#   stop_service()  -> cleanup on stop
#   reload_service() -> reload configuration (SIGHUP)
#   service_triggers() -> procd config change trigger
#
# IMPORTANT: Configure password before first start:
#   uci set acctl.@acctl[0].password='your_secure_password'
#   uci commit acctl

START=95
STOP=05
USE_PROCD=1

NAME_AC=acser
PROG_AC=/usr/bin/acser

# UCI config path
CONFIG_ACCTL=/etc/config/acctl

# Load configuration values
load_config() {
	config_load acctl
	config_get_bool enabled  acctl acctl enabled 1
	config_get interface     acctl acctl interface br-lan
	config_get port          acctl acctl port 7960
	config_get brditv        acctl acctl brditv 30
	config_get reschkitv     acctl acctl reschkitv 300
	config_get msgitv        acctl acctl msgitv 3
	config_get daemon        acctl acctl daemon 1
	config_get_bool debug    acctl acctl debug 0
}

# Build acser command line arguments
build_ac_args() {
	args="-n $interface"
	args="$args -p $port"
	args="$args -c $reschkitv"
	args="$args -m $msgitv"
	args="$args -b $brditv"
	[ "$daemon" = "1" ] && args="$args -d"
	[ "$debug" = "1" ] && args="$args -l 4"
}

# Check prerequisites
check_prereq() {
	[ -x "$PROG_AC" ] || {
		echo "acser not found or not executable at $PROG_AC"
		return 1
	}

	# Ensure database directory exists
	[ -d /etc/acctl ] || mkdir -p /etc/acctl

	# Ensure password is set
	local password
	config_get password acctl acctl password ""
	if [ -z "$password" ]; then
		echo "WARNING: No password configured! Set 'option password' in $CONFIG_ACCTL"
		echo "         Run: uci set acctl.@acctl[0].password='your_secure_password'"
		echo "         Then: uci commit acctl && /etc/init.d/acctl restart"
	fi

	return 0
}

# ---- procd service entry point ----

start_service() {
	config_load acctl
	load_config

	[ "$enabled" != "1" ] && {
		echo "AC Controller is disabled (option enabled='$enabled')"
		return 1
	}

	check_prereq || return 1

	build_ac_args

	procd_open_instance acser
	procd_set_param command "$PROG_AC" $args
	# procd_set_param respawn <threshold> <timeout> <retry_delay>
	#   threshold: max crash count within timeout before giving up
	#   timeout:   window in seconds for counting crashes
	#   retry_delay: seconds to wait before restarting
	procd_set_param respawn 3600 5 10
	procd_set_param file "$CONFIG_ACCTL"
	procd_set_param pidfile /var/run/acser.pid
	procd_close_instance

	echo "AC Controller started (acser, port=$port)"
}

stop_service() {
	echo "Stopping AC Controller..."
	[ -f /var/run/acser.pid ] && {
		kill -TERM "$(cat /var/run/acser.pid)" 2>/dev/null
		rm -f /var/run/acser.pid
	}
	# Also kill by name if pidfile missing
	# FIX: prefer pidfile kill to avoid killing wrong processes
	killall acser 2>/dev/null  # fallback if pidfile missing
	echo "AC Controller stopped"
}

reload_service() {
	echo "Reloading AC Controller configuration..."
	[ -f /var/run/acser.pid ] && {
		kill -HUP "$(cat /var/run/acser.pid)" 2>/dev/null
	}
}

service_triggers() {
	procd_add_config_trigger "config.change" "acctl" "/etc/init.d/acctl" reload
}

# ---- Convenience commands (callable as /etc/init.d/acctl status) ----

status() {
	echo "=== AC Controller Status ==="
	if pgrep -x acser > /dev/null 2>&1; then
		echo "acser:  RUNNING (pid=$(pgrep -x acser))"
	else
		echo "acser:  STOPPED"
	fi
	if pgrep -x apctl > /dev/null 2>&1; then
		echo "apctl:  RUNNING (pid=$(pgrep -x apctl))"
	else
		echo "apctl:  STOPPED"
	fi

	echo ""
	echo "=== Database ==="
	if [ -f /etc/acctl/ac.json ]; then
		echo "JSON database: /etc/acctl/ac.json ($(wc -c < /etc/acctl/ac.json) bytes)"
		# Use acctl-cli for stats if available
		if [ -x /usr/bin/acctl-cli ]; then
			/usr/bin/acctl-cli stats 2>/dev/null || echo "Stats unavailable"
		fi
	else
		echo "Database not found at /etc/acctl/ac.json"
	fi
}
