# luci/applications/luci-app-acctl/root/etc/uci-defaults/luci-app-acctl
#!/bin/sh
# LuCI app-acctl post-install script
# Creates default UCI config if not exists

[ -f /etc/config/acctl ] && exit 0

cat > /etc/config/acctl <<'EOF'
config acctl 'acctl'
	option enabled     '1'
	option interface   'br-lan'
	option port        '7960'
	option password    'acctl@2024'
	option brditv      '30'
	option reschkitv   '300'
	option msgitv      '3'
	option daemon      '1'
	option debug       '0'

config node_default 'default'
	list ssid 'OpenWrt-AC'
	option encryption 'psk2'
	option key 'changeme123'
	option channel_mode '11n'
	option channel 'auto'
	option signal 'auto'
EOF

# Initialize database directory (JSON backend, no sqlite3)
[ -d /etc/acctl ] || mkdir -p /etc/acctl

# Create default database with IP pool
cat > /etc/acctl/ac.json <<'EOF'
{
  "ac": {
    "uuid": "ac-$(cat /proc/sys/kernel/random/uuid)",
    "name": "OpenWrt-AC"
  },
  "resource": {
    "ip_start": "192.168.1.200",
    "ip_end": "192.168.1.254",
    "ip_mask": "255.255.255.0"
  },
  "aps": [],
  "groups": [],
  "alarms": [],
  "firmwares": []
}
EOF

# Enable and start the service
/etc/init.d/acctl enable
/etc/init.d/acctl start

exit 0
